# externals/wbroker-rs.service
[Unit]
Description=wbroker-rs Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/wbroker-rs
ExecStart=~/bin/wbroker-rs
Restart=always
RestartSec=5s

[Install]
WantedBy=default.target
