# externals/wbroker-rs.service
[Unit]
Description=wbroker-rs Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/wbroker-rs
ExecStart=/opt/wbroker-rs/bin/wbroker-rs
Environment=WBROKER_CONFIG=/opt/wbroker-rs/etc/wbroker-rs.toml
Restart=always
RestartSec=5s

[Install]
WantedBy=default.target
